# wxrobot/helperServer/configure/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=~/configure/supervisor.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
username=jerry              ; (default is no username (open server))
password=jinren1987               ; (default is no password (open server))

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9001        ; (ip_address:port specifier, *:port for all iface)
username=jerry              ; (default is no username (open server))
password=jinren1987               ; (default is no password (open server))

[supervisord]
logfile=~/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=~/configure/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY="value"     ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:hub_server]
command=python %(ENV_HOME)s/weixin_helper/hub/hub_server.py
stderr_logfile = %(ENV_HOME)s/log/supervisor/hub_server_stderr.log
stdout_logfile = %(ENV_HOME)s/log/supervisor/hub_server_stdout.log

[program:helper_server]
command=python %(ENV_HOME)s/weixin_helper/helper/hepler_server.py
stderr_logfile = %(ENV_HOME)s/log/supervisor/hepler_server_stderr.log
stdout_logfile = %(ENV_HOME)s/log/supervisor/hepler_server_stdout.log

[program:file_server]
command=python %(ENV_HOME)s/weixin_helper/web/web_server.py
stderr_logfile = %(ENV_HOME)s/log/supervisor/file_server_stderr.log
stdout_logfile = %(ENV_HOME)s/log/supervisor/file_server_stdout.log


;[program:web_main]
;command=python %(ENV_HOME)s/workspace/pythonServer/explore_server/web_admin/web_main.py
;stderr_logfile = %(ENV_HOME)s/logs/supervisor/web_main_stderr.log
;stdout_logfile = %(ENV_HOME)s/logs/supervisor/web_main_stdout.log

;配合 nginx
;[program:web_main]
;command=python %(ENV_HOME)s/workspace/pythonServer/explore_server/web_main.py --port=808%(process_num)d
;process_name = cerium_web_server808%(process_num)d
;stderr_logfile = %(ENV_HOME)s/supervisor/logs/cerium_web_server_stderr.log
;stdout_logfile = %(ENV_HOME)s/supervisor/logs/cerium_web_server_stdout.log
;numprocs = 2
;numprocs_start = 0


[group:server]
programs=hub_server,helper_server

[group:web]
programs=file_server


